# assets/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND -k start
autorestart=true
